# upstart.conf
description "uWSGI server instance configured to serve <project>"

start on runlevel [2345]
stop on runlevel [!2345]

setuid www-data
setgid www-data

chdir /vagrant
exec uwsgi --emperor emperor.ini
